# supervisord.conf
[supervisord]
nodaemon=true

[program:tor]
priority=1
command=/usr/bin/tor

[program:bitcoind]
priority=2
command=/usr/bin/bitcoind

